# rootfs/etc/s6-overlay/s6-rc.d/omeka-s-setup/up
/etc/s6-overlay/scripts/omeka-s-setup.sh

